# 3.4/supervisor.conf
[program:app]
command=uwsgi -x /opt/uwsgi/uwsgi.xml --enable-threads
user=www-data
directory=/opt/www
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=QUIT
